# 04/70-alarm-keyboard.rules
ENV{ID_BUS}=="usb", ENV{ID_INPUT_KEYBOARD}=="1", ACTION=="add", RUN+="/sbin/modprobe keyboard_alarm"
